# pihole.cron
# m h dom mon dow cmd

# Weekly backups
0 1 * * 0 /home/ben/pihole/teleporter.sh &>/dev/null
